# systemd/chicken-gate-web.service
[Unit]
Description=Chicken Gate Web Interface
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sw/chicken-gate
Environment=PATH=/home/pi/.local/bin:/home/pi/sw/chicken-gate/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/pi/sw/chicken-gate/src
ExecStart=/home/pi/.local/bin/chicken-gate-web --port 5000
Restart=always
RestartSec=10s
TimeoutStartSec=30s
TimeoutStopSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
