# systemd/chicken-gate-web.service
[Unit]
Description=Chicken Gate Web Interface
After=multi-user.target network-online.target chicken-gate.service
Wants=network-online.target
Requires=chicken-gate.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sw/chicken-gate
Environment=PATH=/home/pi/sw/chicken-gate/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/pi/sw/chicken-gate/src
ExecStart=/home/pi/sw/chicken-gate/scripts/chicken-gate-web --port5000
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
